# scripts/crawl.crontab
47 * * * * /home/akkartik/readwarp/crawl >> /tmp/readwarp_crawl.out 2>> /tmp/readwarp_crawl.err
