# scripts/crontab
5 3 * * * /home/akkartik/readwarp/scripts/daily_email >> /tmp/daily_email.out 2>> /tmp/daily_email.err
35 * * * * /home/akkartik/readwarp/scripts/backup-data >> /tmp/backup.out 2>> /tmp/backup.err
